# home/dot_config/powershell/profile.d/init.ps1.tmpl
{{- template "powershell/profile.d/init.ps1" . -}}
