# home/dot_config/powershell/profile.d/variables.ps1.tmpl
{{- template "powershell/profile.d/variables.ps1" . -}}
