# home/dot_config/powershell/profile.d/web.ps1.tmpl
{{- template "powershell/profile.d/web.ps1" . -}}
