# home/Documents/PowerShell/profile.d/aliases.ps1.tmpl
{{- template "powershell/profile.d/aliases.ps1" . -}}
