# home/Documents/PowerShell/profile.d/modules.ps1.tmpl
{{- template "powershell/profile.d/modules.ps1" . -}}
