# home/Documents/PowerShell/profile.d/prompt.ps1.tmpl
{{- template "powershell/profile.d/prompt.ps1" . -}}
