# home/Documents/PowerShell/profile.d/readline.ps1.tmpl
{{- template "powershell/profile.d/readline.ps1" . -}}
