# home/Documents/PowerShell/profile.d/utilities.ps1.tmpl
{{- template "powershell/profile.d/utilities.ps1" . -}}
